# crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

*/1 * * * * root sleep $((RANDOM%60)); cd /opt/ews/; /usr/bin/python ews.py -c /data/ews/conf/
